# solaar-redivert.service
[Unit]
Description=Solaar Forward Button redivert watcher
After=graphical-session.target

[Service]
Type=simple
# __INSTALL_DIR__ is replaced by install.sh
ExecStart=/usr/bin/python3 __INSTALL_DIR__/redivert_watch.py
WorkingDirectory=__INSTALL_DIR__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
